# xray-manager.conf
[program:init]
command=/bin/sh -c 'test -f /root/xray-manager/data/geoip.dat || (cd /root/xray-manager/data && curl -fSL https://hub.543083.xyz/https://github.com/Loyalsoldier/v2ray-rules-dat/releases/latest/download/geoip.dat -o geoip.dat && curl -fSL https://hub.543083.xyz/https://github.com/Loyalsoldier/v2ray-rules-dat/releases/latest/download/geosite.dat -o geosite.dat) && ln -sf /root/xray-manager/data/geoip.dat /root/xray-manager/bin/geoip.dat && ln -sf /root/xray-manager/data/geosite.dat /root/xray-manager/bin/geosite.dat && echo "init done" && exit 0'
autostart=true
autorestart=false
startsecs=0
priority=10

[program:xray]
command=/root/xray-manager/bin/xray run -config /root/xray-manager/config/xray-multi-socks.json
autostart=true
autorestart=true
startsecs=2
priority=20
stdout_logfile=/root/xray-manager/logs/xray.log
stderr_logfile=/root/xray-manager/logs/xray.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:xray-panel]
command=python3 /root/xray-manager/app.py --host 0.0.0.0 --port 54321 --xray-config /root/xray-manager/config/xray-multi-socks.json --xray-binary /root/xray-manager/bin/xray --service xray
autostart=true
autorestart=true
startsecs=2
priority=30
stdout_logfile=/root/xray-manager/logs/panel.log
stderr_logfile=/root/xray-manager/logs/panel.err
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3

[group:xray-all]
programs=init,xray,xray-panel
