# xray-manager-v3.conf
[program:init]
command=/bin/sh -c 'mkdir -p /root/xray-manager-v3/{config,data,bin,backup,logs} && test -f /root/xray-manager-v3/data/token || echo 123456 > /root/xray-manager-v3/data/token && chmod 600 /root/xray-manager-v3/data/token 2>/dev/null || true && ln -sf /root/xray-manager-v3/data/geoip.dat /root/xray-manager-v3/bin/geoip.dat && ln -sf /root/xray-manager-v3/data/geosite.dat /root/xray-manager-v3/bin/geosite.dat && echo init done'
autostart=true
autorestart=false
startsecs=0
priority=10

[program:xray]
command=/root/xray-manager-v3/bin/xray run -config /root/xray-manager-v3/config/xray-multi-socks.json
autostart=true
autorestart=true
startsecs=2
priority=20
stdout_logfile=/root/xray-manager-v3/logs/xray.log
stderr_logfile=/root/xray-manager-v3/logs/xray.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:xray-panel]
command=/usr/local/bin/node /root/xray-manager-v3/index.js
directory=/root/xray-manager-v3
autostart=true
autorestart=true
startsecs=2
priority=30
stdout_logfile=/root/xray-manager-v3/logs/panel.log
stderr_logfile=/root/xray-manager-v3/logs/panel.err
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
environment=NODE_ENV="production",PORT="54321"

[group:xray-all]
programs=init,xray,xray-panel
