# xray-manager.service
[Unit]
Description=Xray Manager Web Panel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/xray-manager
ExecStart=/root/xray-manager/.venv/bin/python3 /root/xray-manager/app.py --host 0.0.0.0 --port 54321 --xray-config /root/xray-manager/config/xray-multi-socks.json --xray-binary /root/xray-manager/bin/xray --service xray-multi-socks.service
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/root/xray-manager/logs/panel.log
StandardError=append:/root/xray-manager/logs/panel.err

[Install]
WantedBy=multi-user.target
